# supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true

[program:gunicorn]
command=/usr/local/bin/gunicorn --pythonpath /sfblog sfblog_project.wsgi:application -b 0.0.0.0:8000
stdout_events_enabled=true
stderr_events_enabled=true
